# rootfs/etc/services.d/003-drush-site-install/run
#!/bin/sh
if [ -f /var/www/app/docroot/sites/default/settings.php ]; then
  drush si $DRUPAL_INSTALL_PROFILE -y 
fi